# .config/profile.d/nvm.bash
#!/usr/bin/env bash

export NVM_DIR="$HOME/.nvm"

# shellcheck disable=SC1091
[ -s "/usr/local/opt/nvm/nvm.sh" ] && source "/usr/local/opt/nvm/nvm.sh"
