# .config/profile.d/ps1.bash
#!/usr/bin/env bash

export PS1='\[\033[01;34m\]\w\[\033[0m\] '
