# .config/profile.d/rubygems.bash
#!/usr/bin/env bash

export GEMRC="$HOME/.config/rubygems/gemrc"
