# .config/profile.d/visual-studio-code.bash
#!/usr/bin/env bash

alias update_vscode_extension_list='code --list-extensions > "$HOME/Library/ApplicationSupport/Code/User/extensions.txt"'
